# infra/systemd/otakustream-certbot-renew.timer
[Unit]
Description=Twice-daily TLS renewal for Otakustream

[Timer]
OnCalendar=*-*-* 03,15:00:00
Persistent=true

[Install]
WantedBy=timers.target
